# salt/states/java.service.tmpl
[Unit]
Description={{ name }}
After=syslog.target

[Service]
ExecStart={{ executable }}
User={{ user }}

[Install]
WantedBy=multi-user.target
